# latest/fs/etc/services.d/syslog/run
#!/usr/bin/execlineb -P

exec rsyslogd -n
