# debug/fs/etc/services.d/mailhog/run
#!/usr/bin/execlineb -P

exec /usr/local/bin/mailhog -api-bind-addr=0.0.0.0:80 -ui-bind-addr=0.0.0.0:80
